# spectre_IDS/service/spectre-capture.service
[Unit]
Description=Spectre IDS Packet Capture Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/spectre/spectre_IDS
ExecStart=/opt/spectre/venv/bin/python capture.py
Restart=always
RestartSec=5
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
