# spectre_IDS/service/spectre-ids.target
[Unit]
Description=Spectre IDS Full Stack
Requires=spectre-capture.service spectre-analyzer.service spectre-web.service
After=network.target
