# spectre_IDS/service/spectre-web.service
[Unit]
Description=Spectre IDS Web Interface
After=network.target spectre-analyzer.service

[Service]
Type=simple
WorkingDirectory=/opt/spectre
ExecStart=/opt/spectre/venv/bin/gunicorn -w 1 -b 127.0.0.1:5000 webapp.app:app
Restart=always
RestartSec=5
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target